# sv/dashboard/run
#!/bin/bash

sv start tyk || exit 1

source /etc/envvars
env

cd /tyk-analytics
exec 2>&1
exec ./tyk-analytics
